# deploy/gunicorn-ecomm-api.conf
description "Gunicorn ECOMM API"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5
exec /CHANGE/THE/PATH/ecomm-api/gunicorn.sh